# px4_ws/devel/etc/catkin/profile.d/05.catkin_make.bash
/home/odroid/MAV-Project/px4_ws/devel/.private/catkin/etc/catkin/profile.d/05.catkin_make.bash